# rootfs/etc/inittab
/bin/serialin
/bin/serialout
/bin/output
/bin/shell
